# rc.local

sudo /usr/lib/cgi-bin/start.cgi
screen -dmS motion motion

_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi
